# config/supervisord.conf
[supervisord]
nodaemon=true

[program:auxilium]
command=/app/launch.sh
stopwaitsecs=10
stopsignal=TERM
